# archive/zero-api-scaffold/automation/launchd/com.matt.spiceflow.run-all.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
 <key>Label</key><string>com.matt.spiceflow.run-all</string>
 <key>ProgramArguments</key>
 <array>
   <string>/bin/zsh</string>
   <string>-lc</string>
   <string>cd /Users/matt/src/spiceflow-social && source .venv/bin/activate && SPICEFLOW_ALLOW_RUN=1 python src/run_all.py --horizon-days 45</string>
 </array>
 <key>StartCalendarInterval</key><dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
 <key>RunAtLoad</key><true/>
 <key>StandardOutPath</key><string>/tmp/spiceflow.out.log</string>
 <key>StandardErrorPath</key><string>/tmp/spiceflow.err.log</string>
</dict></plist>
